# files/rc.local.j2
losetup /dev/loop0 /opt/stack/data/stack-volumes-default-backing-file 
losetup /dev/loop1 /opt/stack/data/stack-volumes-lvmdriver-1-backing-file 
exit 0